# quadlets/vaultwarden/vaultwarden.container
[Unit]
Description=Vaultwarden quadlet

[Container]
Image=ghcr.io/dani-garcia/vaultwarden:latest
AutoUpdate=registry

Network=apps.network
NetworkAlias=vaultwarden

# Environment=SIGNUPS_ALLOWED=false
Environment=DOMAIN=https://nemnix.site

Volume=vaultwarden-data.volume:/data:Z

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
