# postgres/postgres.container
[Unit]
Description=PostgreSQL quadlet
Requires=db.network postgres-data.volume
After=db.network postgres-data.volume

[Container]
Image=docker.io/library/postgres:17.7-bookworm
AutoUpdate=registry

EnvironmentFile=./postgres.env

Network=db
NetworkAlias=pgsql

Volume=postgres-data.volume:/var/lib/postgresql/data:Z

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
